# etc/plptools.in
#!/bin/sh
#
# plptools        Starts ncpd/plpfuse.
#
# chkconfig: 2345 45 10
# description: This facility enables connectivity to an EPOC PDA.

# Source function library.
. /etc/rc.d/init.d/functions

[ -f @prefix@/sbin/ncpd ] || exit 0
[ -f @prefix@/sbin/plpfuse ] || exit 0
[ -f @prefix@/sbin/plpprintd ] || exit 0
[ -f /etc/sysconfig/plptools ] || exit 0
. /etc/sysconfig/plptools

start() {
	echo "Starting EPOC support ..."
	RETVAL=0
	if test "$START_NCPD" = "yes" ; then
		echo -n "  ncpd: "
		daemon @prefix@/sbin/ncpd $NCPD_ARGS
		RETVAL=$?
		echo
	fi
        if [ $RETVAL -eq 0 ] ; then
		if test "$START_PLPFUSE" = "yes" ; then
			echo -n "  plpfuse: "
			daemon @prefix@/sbin/plpfuse $PLPFUSE_ARGS
                	RETVAL=$?
			echo
		fi
	fi
	if [ $RETVAL -eq 0 ] ; then
		if test "$START_PLPPRINTD" = "yes" ; then
			echo -n "  plpprintd: "
			daemon @prefix@/sbin/plpprintd $PLPPRINTD_ARGS
			RETVAL=$?
			echo
		fi
	fi
	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/plptools
	return $RETVAL
}

stop() {
	echo "Stopping EPOC support ..."
	RETVAL=0
	if test "$START_PLPPRINTD" = "yes" ; then
		echo -n "  plpprintd: "
		killproc plpprintd
		echo
	fi
	if test "$START_PLPFUSE" = "yes" ; then
		echo -n "  plpfuse: "
		fusermount -u /mnt/epoc
		echo
	fi
	if test "$START_NCPD" = "yes" ; then
		echo -n "  ncpd: "
		killproc ncpd
		RETVAL=$?
		echo
	fi
	[ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/plptools
	return $RETVAL
}

restart() {
	stop
	start
}

# See how we were called.
case "$1" in
  start)
  	start
	;;
  stop)
  	stop
	;;
  status)
	status ncpd
	status plpprintd
	;;
  restart|reload)
  	restart
	;;
  condrestart)
  	test -f /var/lock/subsys/plptools && restart || :
	;;
  *)
	echo "Usage: plptools {start|stop|status|restart|reload|condrestart}"
	exit 1
esac

exit $?

