# etc/psion.in
#!/bin/sh
#
# psion        Starts ncpd/plpnfsd.
#
# chkconfig: 2345 45 70
# description: This facility enables connectivity to a Psion series 5.

# Source function library.
. /etc/rc.d/init.d/functions

[ -f @prefix@/sbin/ncpd ] || exit 0
[ -f @prefix@/sbin/plpnfsd ] || exit 0

MGETTY_HASPLP=false

if grep -qs ^/PLP/ /etc/mgetty+sendfax/login.config ; then
	mrun=`pidofproc mgetty`
	test -n "$mrun" && MGETTY_HASPLP=true
fi

start() {
  	if $MGETTY_HASPLP ; then
		echo "NOT Starting ncpd because mgetty configured for PLP"
		exit 0
	fi
	echo -n "Starting Psion support (ncpd): "
	daemon @prefix@/sbin/ncpd
	RETVAL=$?
	echo
	if [ $RETVAL -eq 0 ] ; then
		echo -n "Starting Psion support (plpnfsd): "
		daemon @prefix@/sbin/plpnfsd
		RETVAL=$?
		echo
	fi
	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/psion
	return $RETVAL
}

stop() {
	echo -n "Stopping Psion support: "
	killproc plpnfsd -HUP
	while true ; do
		test -z "`pidofproc plpnfsd`" && break;
		sleep 1 # allow plpnfsd flushing it's cache
	done
	killproc ncpd
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/psion
	return $RETVAL
}

restart() {
	stop
	start
}

# See how we were called.
case "$1" in
  start)
  	start
	;;
  stop)
  	stop
	;;
  status)
	status plpnfsd
	status ncpd
	;;
  restart|reload)
  	restart
	;;
  condrestart)
  	test -f /var/lock/subsys/psion && restart || :
	;;
  *)
	echo "Usage: psion {start|stop|status|restart|reload|condrestart}"
	exit 1
esac

exit $?

