# etc/psion.SuSE.in
#!/bin/sh
#
### BEGIN INIT INFO
# Provides:       psion
# Required-Start: $syslog
# Required-Stop:
# Default-Start:  3 5
# Default-Stop:   0 1 2 6
# Description:    Start Psion support
### END INIT INFO

# Source function library.
. /etc/rc.status

[ -f @prefix@/sbin/ncpd ] || exit 0
[ -f @prefix@/sbin/plpprintd ] || exit 0
[ -f /etc/sysconfig/psion ] || exit 0
. /etc/sysconfig/psion

rc_reset

start() {
	echo "Starting Psion support ..."
	if [ "$START_NCPD" = "yes" ] ; then
		echo -n "  ncpd: "
		startproc -f @prefix@/sbin/ncpd $NCPD_ARGS
		rc_status -v
	fi
	if [ "$START_PLPPRINTD" = "yes" ] ; then
		echo -n "  plpprintd: "
		startproc -f @prefix@/sbin/plpprintd $PLPPRINTD_ARGS
		rc_status -v
	fi
	rc_status && touch /var/lock/subsys/psion
}

stop() {
	echo "Stopping Psion support ..."
	if [ "$START_PLPPRINTD" = "yes" ] ; then
		echo -n "  plpprintd: "
		killproc -TERM @prefix@/sbin/plpprintd
		rc_status -v
	fi
	if [ "$START_NCPD" = "yes" ] ; then
		echo -n "  ncpd: "
		killproc -TERM @prefix@/sbin/ncpd
		rc_status -v
	fi
	rc_status && rm -f /var/lock/subsys/psion
}

restart() {
	stop
	start
}

# See how we were called.
case "$1" in
  start)
  	start
	;;
  stop)
  	stop
	;;
  status)
	echo -n "Checking for service psion: "
	checkproc @prefix@/sbin/plpprintd && \
		checkproc @prefix@/sbin/ncpd
	rc_status -v
	;;
  restart|reload)
  	restart
	;;
  condrestart)
  	test -f /var/lock/subsys/psion && restart || :
	;;
  *)
	echo "Usage: psion {start|stop|status|restart|reload|condrestart}"
	exit 1
esac

rc_exit
